# .deployment/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; Messenger async consumer - processes queued jobs (emails, notifications, etc.)
[program:messenger-async]
command=php /app/bin/console messenger:consume async --time-limit=3600 --memory-limit=256M --sleep=1 -vv
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stopsignal=SIGTERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data

; Scheduler consumer - runs scheduled/periodic tasks
[program:scheduler]
command=php /app/bin/console messenger:consume scheduler_default --time-limit=3600 --memory-limit=128M -vv
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stopsignal=SIGTERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data
